# deploy/systemd/security-daily-ban-digest.service
[Unit]
Description=Security Daily Ban Digest
After=network-online.target fail2ban.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/security-report-bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/default/security-report-bot
ExecStart=/opt/security-report-bot/.venv/bin/python -m app.daily_digest
