# deploy/systemd/security-manual-denylist-sync.service
[Unit]
Description=Sync Security Manual Denylist to UFW
After=network-online.target ufw.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/security-report-bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/default/security-report-bot
ExecStart=/opt/security-report-bot/.venv/bin/python -m app.manual_denylist sync
