# deploy/systemd/security-manual-denylist-sync.path
[Unit]
Description=Watch Security Manual Denylist file

[Path]
PathModified=/etc/security-report-bot/manual-denylist.txt
Unit=security-manual-denylist-sync.service

[Install]
WantedBy=multi-user.target
