# deploy/systemd/security-allowlist-sync.path
[Unit]
Description=Watch security-report-bot scan allowlist file

[Path]
PathModified=/etc/security-report-bot/scan-whitelist.txt
Unit=security-allowlist-sync.service

[Install]
WantedBy=multi-user.target
