# deploy/systemd/security-fail2ban-ban-snapshot.service
[Unit]
Description=Snapshot Active Fail2Ban Web Bans
After=network-online.target fail2ban.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/security-report-bot
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/default/security-report-bot
ExecStart=/opt/security-report-bot/.venv/bin/python -m app.fail2ban_persistent snapshot
